# services/personal-secretary/personal-secretary.service
[Unit]
Description=personal-secretary Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=theau
Group=theau
WorkingDirectory=/var/lib/personal-secretary
EnvironmentFile=/etc/personal-secretary/secrets.env
Environment=PERSONAL_SECRETARY_DATA_DIR=/var/lib/personal-secretary
Environment=TZ=Europe/Brussels
ExecStart=/var/lib/personal-secretary/venv/bin/python /var/lib/personal-secretary/main.py
Restart=on-failure
RestartSec=10
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/personal-secretary

[Install]
WantedBy=multi-user.target
